# crontab
MAILTO=<your email address>
30 1 * * *	cd $HOME/globus-sync-example && BASH_ENV=/etc/profile ./run_globus_sync.sh
